# usr/lib/systemd/system/droidcam2v4l2.service.d/10-module.conf
[Unit]
ConditionPathExists=
ConditionPathExists=

[Service]
ExecStartPre=-/usr/sbin/insmod /vendor/lib/modules/v4l2loopback.ko
