# usr/lib/systemd/system/android-mount.service.d/40-slot-b.conf
[Service]
ExecStartPre=-/bin/bash -c 'if [ "$(grep -o 'androidboot\.slot_suffix=..' /proc/cmdline |  cut -d "=" -f2)" == "_a" ]; then ln -s /dev/disk/by-partlabel/boot_a /dev/disk/by-partlabel/boot; else ln -s /dev/disk/by-partlabel/boot_b /dev/disk/by-partlabel/boot; fi'
ExecStartPre=-/bin/bash -c 'if [ "$(grep -o 'androidboot\.slot_suffix=..' /proc/cmdline |  cut -d "=" -f2)" == "_a" ]; then ln -s /dev/disk/by-partlabel/dtbo_a /dev/disk/by-partlabel/dtbo; else ln -s /dev/disk/by-partlabel/dtbo_b /dev/disk/by-partlabel/dtbo; fi'
