# usr/lib/systemd/system/android-mount.service.d/50-geoclue.conf
[Service]
ExecStartPost=-/usr/bin/mkdir -p /etc/geoclue/conf.d/
ExecStartPost=-/usr/bin/chown 32011:32011 /etc/geoclue/conf.d/ -R
