# usr/lib/systemd/system/android-mount.service.d/60-modules.conf
[Service]
ExecStartPre=-/usr/libexec/module2vendor
