# usr/lib/systemd/user/sm.puri.Phosh.service.d/10-120hz.conf
[Service]
ExecStartPost=-/usr/libexec/krypton-120hz
